# deploy/systemd/orca-ingest-live.service
[Unit]
Description=Orca near-live ingest service
After=network-online.target docker.service
Wants=network-online.target docker.service

[Service]
Type=simple
User=lynchej
WorkingDirectory=/home/lynchej/orca_polymarkets
EnvironmentFile=/etc/orca.env
ExecStart=/bin/bash /home/lynchej/orca_polymarkets/scripts/run_ingest_live_vm.sh
Restart=always
RestartSec=15
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
